# platform/bin/fusion
##!/bin/sh
# chkconfig: 2345 99 50
# description: fusion management script
#
# Adapt the following lines to your configuration
FUSION_HOME=/usr/local/fusion
JAVA_HOME=/usr/java/latest/jre
PID_FILE=/var/run/jsvc-fusion.pid

FUSION_OPTS="-Xms512m -Xmx8192m"

CLASSPATH=.:server.jar

cd $FUSION_HOME

case "$1" in
  start)
    /usr/bin/jsvc \
    -home $JAVA_HOME \
    -jvm server \
    -cwd $FUSION_HOME \
    -pidfile $PID_FILE \
    -outfile $FUSION_HOME/logs/stdout.log \
    -errfile $FUSION_HOME/logs/stderr.log \
    $FUSION_OPTS \
    -cp $CLASSPATH \
    lsfusion.server.logics.BusinessLogicsBootstrap
    exit $?
    ;;

  stop)
    /usr/bin/jsvc \
    -home $JAVA_HOME \
    -stop \
    -pidfile $PID_FILE \
    lsfusion.server.logics.BusinessLogicsBootstrap
    exit $?
    ;;

  *)
    echo "Usage fusion.sh start/stop"
    exit 1;;
esac
